# little.sh
#!/bin/sh

# æ£€æŸ¥è¾“å…¥å‚æ•°
if [ $# -eq 0 ]; then
    echo "ç”¨æ³•: $0 \"è¦å¼€æœºå¯åŠ¨çš„å®Œæ•´å‘½ä»¤\""
    echo "ç¤ºä¾‹: $0 
    exit 1
fi

COMMAND="$1"

# æ£€æŸ¥æ˜¯å¦ä¸ºrootç”¨æˆ·
if [ "$(id -u)" != "0" ]; then
    echo "é”™è¯¯: æ­¤è„šæœ¬éœ€è¦rootæƒé™è¿è¡Œ"
    exit 1
fi

# æ£€æŸ¥Alpineç³»ç»Ÿ
if ! grep -q "Alpine" /etc/os-release; then
    echo "é”™è¯¯: æ­¤è„šæœ¬ä»…é€‚ç”¨äºŽAlpineç³»ç»Ÿ"
    exit 1
fi

# åˆ›å»ºæœ¬åœ°å¯åŠ¨è„šæœ¬
LOCAL_SCRIPT="/usr/local/bin/startup_command.sh"

cat > "$LOCAL_SCRIPT" << EOF
#!/bin/sh
# å¼€æœºå¯åŠ¨è„šæœ¬
# ç”±ä¸€é”®è„šæœ¬ç”ŸæˆäºŽ $(date)

# ç­‰å¾…ç½‘ç»œå°±ç»ª
sleep 10

# æ‰§è¡Œç”¨æˆ·å‘½ä»¤
$COMMAND
EOF

chmod +x "$LOCAL_SCRIPT"

# åˆ›å»ºOpenRCæœåŠ¡æ–‡ä»¶
SERVICE_FILE="/etc/init.d/startup_command"

cat > "$SERVICE_FILE" << EOF
#!/sbin/openrc-run
# å¼€æœºå¯åŠ¨æœåŠ¡
# ç”±ä¸€é”®è„šæœ¬ç”ŸæˆäºŽ $(date)

name="startup_command"
description="ç”¨æˆ·è‡ªå®šä¹‰å¼€æœºå¯åŠ¨å‘½ä»¤"

command="$LOCAL_SCRIPT"
command_background=true
pidfile="/var/run/\${name}.pid"

depend() {
    need net
    after firewall
}
EOF

chmod +x "$SERVICE_FILE"

# æ·»åŠ æœåŠ¡åˆ°å¯åŠ¨é¡¹
rc-update add startup_command default

echo "âœ… å¼€æœºå¯åŠ¨å·²è®¾ç½®å®Œæˆ!"
echo "ðŸ“ æœåŠ¡åç§°: startup_command"
echo "ðŸ“‹ å¯åŠ¨å‘½ä»¤: $COMMAND"
echo "ðŸ“ æœ¬åœ°è„šæœ¬: $LOCAL_SCRIPT"
echo "ðŸ”§ æœåŠ¡æ–‡ä»¶: $SERVICE_FILE"
echo ""
echo "é‡å¯åŽç”Ÿæ•ˆï¼Œå¦‚éœ€ç«‹å³æµ‹è¯•å¯è¿è¡Œ: $LOCAL_SCRIPT"
echo "å¦‚éœ€ç§»é™¤æœåŠ¡ï¼Œè¯·è¿è¡Œ: rc-update del startup_command"
